# session/xdg-desktop-portal-villode.service
[Unit]
Description=Villode FileChooser portal backend
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.villode
ExecStart=%h/.local/bin/xdg-desktop-portal-villode
Restart=on-failure
RestartSec=1

[Install]
WantedBy=graphical-session.target
